# ansible/roles/udp2raw/templates/udp2raw-client.service.j2
[Unit]
Description=udp2raw tunnel client
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/udp2raw \
  -c \
  -l {{ vpn_entry_ip }}:{{ wireguard_port }} \
  -r {{ vpn_exit_ip }}:{{ wireguard_port }} \
  -k "{{ udp2raw_password }}" \
  --raw-mode faketcp
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target